# images/probedock-app/config/services.d/unicorn/run
#!/usr/bin/with-contenv bash

rm -f /usr/src/app/tmp/pids/server.pid

cd /usr/src/app

exec rails server -e production -b 0.0.0.0