# images/probedock-app/config/services.d/serf/run
#!/usr/bin/with-contenv bash

pid=$$

sigterm() {
  kill -INT $pid
}

trap sigterm SIGTERM

/opt/bin/serf agent --config-file /etc/serf.conf &

pid=$!

wait