# images/probedock-job/config/services.d/resque/run
#!/usr/bin/with-contenv bash

cd /usr/src/app

exec rake resque:work